# raspberry4/services/micro-ros.service
[Unit]
Description=Micro-ROS service

[Service]
User=USER
ExecStart=/bin/bash --login -c 'source /home/USER/.bashrc; /home/USER/services/micro-ros.sh'
KillSignal=2
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
